# srv.service
[Unit]
Description=SafeSite AI - Construction Site Security Platform
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/chantier-securite
ExecStart=/home/exedev/chantier-securite/safesite -listen :8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
